# root/etc/services.d/wireguard/run
#!/usr/bin/with-contenv bash

_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down ${INTERFACENAME:-wg0}
}

trap _term SIGTERM

wg-quick up ${INTERFACENAME:-wg0}

sleep infinity &

wait
